# deploy/gateway/heyskip-cron.service
[Unit]
Description=Hey Skip Cron Runner
After=network.target

[Service]
Type=oneshot
User=admin
ExecStart=/home/admin/heyskip-agents/deploy/cron-runner.sh
Environment=SKIP_API_URL=${SKIP_API_URL}
Environment=SKIP_API_TOKEN=${SKIP_API_TOKEN}
